# services/lobster-structural-hygiene-audit.timer
[Unit]
Description=Annual structural hygiene audit — deep system structure review
# LOBSTER-MANAGED

[Timer]
OnCalendar=*-03-31 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
