# services/lobster-structural-hygiene-audit.service
[Unit]
Description=Annual structural hygiene audit — deep system structure review
# LOBSTER-MANAGED

[Service]
Type=oneshot
User=lobster
ExecStart=/home/lobster/lobster/scheduled-tasks/dispatch-job.sh structural-hygiene-audit
