# services/lobster-whatsapp-bridge.service
[Unit]
Description=Lobster WhatsApp Bridge (Baileys)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/lobster-workspace/projects/lobster/connectors/whatsapp
Environment=PATH=/home/admin/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=LOBSTER_MESSAGES=/home/admin/messages
Environment=LOBSTER_WORKSPACE=/home/admin/lobster-workspace
Environment=WA_COMMANDS_DIR=/home/admin/messages/wa-commands
Environment=WA_EVENTS_DIR=/home/admin/messages/wa-events
Environment=WA_HEARTBEAT_FILE=/home/admin/lobster-workspace/logs/whatsapp-heartbeat
# User secrets — not in git. See connectors/whatsapp/config.example.env
EnvironmentFile=-/home/admin/.config/lobster/whatsapp.env
# Legacy location also supported
EnvironmentFile=-/home/admin/lobster-workspace/config/whatsapp.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
StandardOutput=append:/home/admin/lobster-workspace/logs/whatsapp-bridge.log
StandardError=append:/home/admin/lobster-workspace/logs/whatsapp-bridge.log
SyslogIdentifier=lobster-whatsapp-bridge

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
