# services/lobster-github-issue-cultivator.timer
[Unit]
Description=Daily GitHub issue cultivator — triage and enrich open issues
# LOBSTER-MANAGED

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
